# config/supervisornode/startup.ini
[program:startupnode]
directory=/web/go-gin-gee
command=/web/go-gin-gee/dist/startupnode-linux-amd64
autostart=true
autorestart=true
stderr_logfile=/web/go-gin-gee/log/startupnode.err
stdout_logfile=/web/go-gin-gee/log/startupnode.log
