# config/supervisorb/startupb.ini
[program:startup]
command=/Users/mazey/Web/Mazey/go-gin-gee/dist/startup-mac-darwin-amd64
autostart=true
autorestart=true
