# config/supervisorb/startupb.ini
[program:startup]
directory=/Users/mazey/Web/Mazey/go-gin-gee
command=/Users/mazey/Web/Mazey/go-gin-gee/dist/startup-mac-darwin-amd64
autostart=true
autorestart=true
stderr_logfile=/Users/mazey/Web/Mazey/go-gin-gee/log/startup.err
stdout_logfile=/Users/mazey/Web/Mazey/go-gin-gee/log/startup.log

[program:shadowsocks2]
directory=/Users/mazey/Web/Mazey/go-gin-gee
command=/Users/mazey/Web/Mazey/go-gin-gee/dist/shadowsocks2-macos-amd64 -s "ss://CHACHA20-IETF-POLY1305:your-password@:8488" -verbose
autostart=true
autorestart=true
stderr_logfile=/Users/mazey/Web/Mazey/go-gin-gee/log/startup.err
stdout_logfile=/Users/mazey/Web/Mazey/go-gin-gee/log/startup.log
