# config/supervisorb/startupb.ini
[program:startup]
directory=/Users/mazey/Web/Mazey/go-gin-gee
command=./dist/startup-mac-darwin-amd64
autostart=true
autorestart=true
