# config/supervisorb/startupb.ini
[program:startup]
directory=/Users/mazey/Web/Mazey/conf/config/supervisorb
command=/Users/mazey/Web/Mazey/conf/config/supervisorb/startup-mac-darwin-amd64.secret.exe
autostart=true
autorestart=true
