# config/supervisorb/startupb.ini
[program:startup]
directory=/Users/mazey/Web/Mazey/go-gin-gee
command=/Users/mazey/Web/Mazey/go-gin-gee/dist/startup-mac-darwin-amd64
autostart=true
autorestart=true
stderr_logfile=/Users/mazey/Web/Mazey/go-gin-gee/log/startup.err
stdout_logfile=/Users/mazey/Web/Mazey/go-gin-gee/log/startup.log
