# config/supervisorjapan/startup.ini
[program:startupjapan]
directory=/web/go-gin-gee
command=/web/go-gin-gee/dist/startupjapan-linux-amd64 -projectpath "/web/go-gin-gee"
autostart=true
autorestart=true
stderr_logfile=/web/go-gin-gee/log/startup.err
stdout_logfile=/web/go-gin-gee/log/startup.log
